# elm-misc/2.3-mobile/HVGA/profile.d/ecore.sh
export ECORE_INPUT_FIX=1
export ECORE_INPUT_TIMEOUT_FIX=0
