# mmfw-sysconf/mmfw-sysconf-e3250/etc/profile.d/mmf.sh
export GST_DEBUG=2
export GST_REGISTRY=/home/app/.gstreamer-0.10/registry.bin
export PULSE_CLIENTCONFIG=/etc/pulse/client.conf
