# elm-misc/2.3-wearable/b3/profile.d/edje.sh
export MULTISENSE_SND_PLAYER="tizen_snd_player"
