# elm-misc/2.3-wearable/b3/profile.d/efl_abort.sh
#export ECORE_ERROR_ABORT=1
#export EINA_LOG_ABORT=1
#export ELM_ERROR_ABORT=1
#export EVAS_DEBUG_ABORT=1
