# elm-misc/2.3-wearable-circle/c1/profile.d/evas.sh
export EVAS_GL_ATLAS_ALLOC_SIZE=128
export EVAS_FONT_DPI=72
export EVAS_GL_NO_BLACKLIST=1
