# elm-misc/2.3-mobile/WVGA/profile.d/eina.sh
export EINA_LOG_DLOG_ENABLE=1
export EINA_LOG_SYSLOG_ENABLE=0
export EINA_LOG_LEVEL=1
export EINA_LOG_LEVELS=ecore_x:4,evas_main:1
